# supervisord.conf
[supervisord]
nodaemon=true

[program:mopidy]
command=/usr/bin/mopidy --config /etc/mopidy/mopidy.conf

[program:snapcast]
command=/usr/bin/snapserver